# system/etc/systemd/system/working-hours-inhibit.timer
[Unit]
Description=Start working hours sleep inhibitor at 08:00

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
Unit=working-hours-inhibit.service

[Install]
WantedBy=timers.target
